# harvia-homekit.service
[Unit]
Description=Harvia Sauna HomeKit Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/opt/harvia-homekit
ExecStart=/usr/bin/python3 /opt/harvia-homekit/main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 